# overlay/etc/init.d/ccapply
#!/sbin/openrc-run

depend() {
    want network-online
}

name="ccapply"
command="/macula/system/macula/current/maculaos"
command_args="config"