# overlay/etc/init.d/macula-firstboot
#!/sbin/openrc-run
# MaculaOS First-Boot Setup Service
# Just prints a message - actual wizard runs after login via profile script

name="macula-firstboot"
description="MaculaOS First-Boot Setup"

CONFIGURED_FLAG="/var/lib/maculaos/.configured"

depend() {
    need net
    after connman dbus avahi-daemon
    before sshd
}

start() {
    # Ensure macula directory exists
    checkpath --directory --owner root:root --mode 0755 /var/lib/maculaos

    # Check if already configured
    if [ -f "${CONFIGURED_FLAG}" ]; then
        einfo "System already configured"
        return 0
    fi

    einfo "First-boot setup required"
    einfo "The setup wizard will run after you log in"
    einfo "Login with: macula / macula"
    return 0
}

stop() {
    return 0
}

status() {
    if [ -f "${CONFIGURED_FLAG}" ]; then
        einfo "System is configured"
        return 0
    fi
    einfo "System not yet configured - login to run wizard"
    return 0
}
